# .docker/app/etc/service/swoole/run
#!/bin/bash

exec 2>&1
source /env

exec /var/www/app/vendor/bin/swoole --env=$APP_ENV
